# deploy/systemd/animeattire.service
[Unit]
Description=AnimeAttire production stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/animeattire
EnvironmentFile=/etc/animeattire/production.env
ExecStart=/usr/bin/docker compose -f docker-compose.vps.yml up -d --build
ExecStop=/usr/bin/docker compose -f docker-compose.vps.yml down
ExecReload=/usr/bin/docker compose -f docker-compose.vps.yml up -d --build
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
